# supervisord.conf

[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "/usr/sbin/apachectl -D FOREGROUND"

